# fake_sqs/docker/service/fake_sqs/run
#!/bin/sh

exec /usr/local/bin/fake_sqs
